# .config/systemd/user/ssl-check.timer
[Unit]
Description=Check SSL certificates

[Timer]
OnBootSec=15min
OnUnitActiveSec=120min

[Install]
WantedBy=timers.target
