# utils/taskwarrior-reporter/.config/systemd/user/taskwarrior-reporter-weekly.timer
[Unit]
Description=Email taskwarrior tasks for a current week

[Timer]
OnBootSec=1min
OnUnitActiveSec=1d
OnCalendar=Sat *-*-*

[Install]
WantedBy=timers.target

