# systemd/.config/systemd/user/diff-notes.timer
[Unit]
Description=Run diff-notes daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

