# systemd/.config/systemd/user/diff-notes.service
[Unit]
Description=diff-notes

[Service]
ExecStart=/home/jubnzv/.local/bin/diff-notes.sh
Type=oneshot

[Install]
WantedBy=default.target
