# systemd/.config/systemd/user/taskwarrior-has-tasks-notification.service
[Unit]
Description=Send notificaiton with number of remaining tasks for today.

[Service]
Environment="DISPLAY=:0"
ExecStart=/home/jubnzv/.local/bin/taskwarrior-notify-has-tasks
Type=oneshot

[Install]
WantedBy=default.target
