# systemd/.config/systemd/user/taskwarrior-end-day.timer
[Unit]
Description=Force taskwarrior-end-day routine at midnight.

[Timer]
OnCalendar=23:55:00

[Install]
WantedBy=timers.target

