# not-in-use/systemd/user/diff-anki.service
[Unit]
Description=diff-anki

[Service]
ExecStart=/home/jubnzv/.local/bin/diff-anki.py
Type=oneshot

[Install]
WantedBy=default.target
