# not-in-use/systemd/user/taskwarrior-active-task-notification.timer
[Unit]
Description=Send notificaiton title of active task.

[Timer]
OnCalendar=*:0/20
Persistent=false

[Install]
WantedBy=timers.target

